# crontab
*/10 * * * * python /usr/src/app/tweet-insta-stories.py >> /var/log/cron.log 2>&1
